# server/livejack-reload.timer
[Unit]
Description=Reload leaking livejack daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
