# server/livejack-reload.service
[Unit]
Description=Restart livejack service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart "livejack@*"

