# broker/livejack-broker.target
[Unit]
Description=Control %p

[Install]
WantedBy=default.target
